# scheduler/systemd/hermes-continuity-archive-dry-run.timer
[Unit]
Description=Hermes Continuity archive dry-run timer example
[Timer]
OnCalendar=weekly
Persistent=false
[Install]
WantedBy=timers.target
