# scheduler/systemd/hermes-continuity-audit.timer
[Unit]
Description=Hermes Continuity report-only audit timer example
[Timer]
OnCalendar=daily
Persistent=false
[Install]
WantedBy=timers.target
