# scheduler/systemd/hermes-continuity-audit.service
[Unit]
Description=Hermes Continuity report-only audit example
[Service]
Type=oneshot
ExecStart=<HCK_RUNTIME_HOME>/bin/hck-memory-audit-dry-run
